# root/etc/services.d/jobber/run
#!/bin/sh

exec with-contenv /jobber/jobberrunner /borgmatic/.jobber
