# root/etc/services.d/jobber/run
#!/usr/bin/sh

/jobber/jobberrunner /borgmatic/.jobber
