# root/etc/services.d/jobber/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh
/jobber/jobberrunner /borgmatic/.jobber
