# root/etc/services.d/borgmatic-logs/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh
exec tail -Fq /var/log/borgmatic/current
