# root/etc/services.d/borgmatic-logs/run
#!/bin/sh

exec with-contenv tail -Fq /var/log/borgmatic/current
