# root/etc/services.d/borgmatic/run
#!/usr/bin/sh
with-contenv

exec snooze -v -H/1 borgmatic -c /borgmatic/config.yaml
