# root/etc/services.d/borgmatic/run
#!/bin/sh

exec with-contenv snooze -v -H/1 borgmatic -c /borgmatic/config.yaml
