# root/etc/services.d/stdout-log/run
#!/usr/bin/with-contenv sh
tail -f /var/log/jobber-runs
